# configs_misc/linux/etc/udev/rules.d/00-usb-keyboards.rules
ACTION=="add", ATTRS{idVendor}=="29ea", ATTRS{idProduct}=="0102", RUN+="/home/jon/bin/kbd_udev", OWNER="jon"
